# dist/linux/byedpi.service
[Unit]
Description=byedpi
Documentation=https://github.com/hufrea/byedpi

[Service]
EnvironmentFile=/etc/default/byedpi
ExecStart=ciadpi $BYEDPI_OPTS
TimeoutStopSec=5s
LimitNOFILE=1048576
LimitNPROC=512
PrivateTmp=true
ProtectSystem=full

[Install]
WantedBy=default.target